# server-setup/configs/9router/9router.service
[Unit]
Description=9router AI Router Service
After=network.target

[Service]
Type=simple
User=__USER__
ExecStart=__NPM_GLOBAL_BIN__/9router --port 8317 --host 0.0.0.0 --no-browser --skip-update
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
